# roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd
Documentation=https://github.com/coreos

[Service]
Type=notify
ExecStart=/usr/local/bin/etcd  --name {{ inventory_hostname }} --cert-file=/etc/etcd/kubernetes.pem --key-file=/etc/etcd/kubernetes-key.pem --peer-cert-file=/etc/etcd/kubernetes.pem --peer-key-file=/etc/etcd/kubernetes-key.pem --trusted-ca-file=/etc/etcd/ca.pem --peer-trusted-ca-file=/etc/etcd/ca.pem --peer-client-cert-auth --client-cert-auth --initial-advertise-peer-urls https://{{ hostvars[inventory_hostname].ansible_all_ipv4_addresses|join('') }}:2380 --listen-peer-urls https://{{ hostvars[inventory_hostname].ansible_all_ipv4_addresses|join('') }}:2380 --listen-client-urls https://{{ hostvars[inventory_hostname].ansible_all_ipv4_addresses|join('') }}:2379,https://127.0.0.1:2379 --advertise-client-urls https://{{ hostvars[inventory_hostname].ansible_all_ipv4_addresses|join('') }}:2379 --initial-cluster-token etcd-cluster-0  --initial-cluster {{ controller_nodes[0] }}=https://{{ controller_private_addresses[0] }}:2380,{{ controller_nodes[1] }}=https://{{ controller_private_addresses[1] }}:2380,{{ controller_nodes[2] }}=https://{{ controller_private_addresses[2] }}:2380 --initial-cluster-state new  --data-dir=/var/lib/etcd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
