# roles/install_vault/templates/vault.service.j2
[Unit]
Description=Vault server 
After=network.target

[Service]
Type=simple
ExecStart=/bin/vault server -config={{ vault_config }}
TimeoutStartSec=0


[Install]
WantedBy=multi-user.target